# package/shairport-sync-airplay2/files/shairport-sync.init
#!/bin/sh /etc/rc.common
START=90
STOP=10

start_service() {
    echo "Starting shairport-sync..."
    /usr/bin/shairport-sync -c /etc/config/shairport-sync -d
}

stop_service() {
    killall shairport-sync || true
}